# luci-app-cussubconverter/root/etc/uci-defaults/40_luci-cussubconverter
#!/bin/sh

#chmod +x /www/cgi-bin/cussubconverter
exit 0
